# caddy/caddy-http.socket
[Unit]
Description=Caddy HTTP Socket

[Socket]
ListenStream=80
Service=caddy.service

[Install]
WantedBy=sockets.target