# caddy/caddy-http.socket
[Unit]
Description=Caddy HTTP Socket

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target