# ufw/wireguard-iptables.service
[Unit]
Description=Host-level WireGuard NAT & Forwarding
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/iptables -t nat -A POSTROUTING -s 10.8.0.0/24 -o eth0 -j MASQUERADE
ExecStart=/sbin/iptables -A FORWARD -i wg0 -o eth0 -j ACCEPT
ExecStart=/sbin/iptables -A FORWARD -i eth0 -o wg0 -m conntrack --ctstate RELATED,ESTABLISHED -j ACCEPT
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target