# ansible/roles/backup/files/backup.service
[Unit]
Description=Vaultwarden Backup Service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-vaultwarden.sh
