# forgejo/forgejo.container
[Unit]
Description=Forgejo via Rootless Sockets (VCS)
Wants=network-online.target
After=network-online.target
Restart=always

[Container]
ContainerName=forgejo
Environment=USER_UID=1001
Environment=USER_GID=1001
Image=codeberg.org/forgejo/forgejo:10
Network=none
Socket=forgejo-ssh.socket
Socket=forgejo-ui.socket
Volume=forgejo-data:/data

[Service]
Restart=always

[Install]
WantedBy=default.target