# forgejo/forgejo-ssh.socket
[Unit]
Description=Forgejo UI HTTP Socket

[Socket]
ListenStream=3.15.122.25:222

[Install]
WantedBy=sockets.target