# forgejo/forgejo.container
[Unit]
Description=Forgejo via Rootless Sockets (VCS)
Requires=forgejo.socket ssh.service
After=forgejo.socket ssh.service

[Container]
ContainerName=forgejo
Environment=USER_UID=1001
Environment=USER_GID=1001
Image=codeberg.org/forgejo/forgejo:10
Network=none
Volume=forgejo-data:/data

Environment=DISABLE_SSH=true
Environment=START_SSH_SERVER=false
Environment=SSH_LISTEN_PORT=222

[Service]
Restart=always

[Install]
WantedBy=default.target