# forgejo/forgejo.container
[Unit]
Description=Forgejo via Rootless Sockets (VCS)
Requires=forgejo.socket
After=forgejo.socket forgejo.socket

[Container]
ContainerName=forgejo
Environment=USER_UID=1001
Environment=USER_GID=1001
Image=localhost/forgejo/forgejo-custom:latest
Network=none
Volume=forgejo-data:/data

Environment=START_SSH_SERVER=true
Environment=SSH_LISTEN_PORT=22

[Service]
Restart=always

[Install]
WantedBy=default.target