# forgejo/forgejo-ssh.socket
[Unit]
Description=Forgejo UI HTTP Socket

[Socket]
ListenStream=222

[Install]
WantedBy=sockets.target