# forgejo/forgejo.container
[Unit]
Description=Forgejo via Rootless Sockets (VCS)
Requires=forgejo.socket ssh.service
After=forgejo.socket ssh.service

[Container]
ContainerName=forgejo
Environment=USER_UID=1001
Environment=USER_GID=1001
Image=codeberg.org/forgejo/forgejo:10
Network=none
Volume=forgejo-data:/data

Environment=START_SSH_SERVER=true
Environment=SSH_PORT=2222
Environment=SSH_LISTEN_PORT=222

[Service]
Restart=always

[Install]
WantedBy=default.target