# forgejo/forgejo.container
[Unit]
Description=Forgejo via Rootless Sockets (VCS)
Requires=forgejo-ui.socket forgejo-ssh.socket
After=forgejo-ui.socket forgejo-ssh.socket
Restart=always
SocketActivated=yes

[Container]
ContainerName=forgejo
Environment=USER_UID=1001
Environment=USER_GID=1001
Image=codeberg.org/forgejo/forgejo:10
Network=none
Volume=forgejo-data:/data

[Service]
Restart=always

[Install]
WantedBy=default.target