# ansible/roles/backup/templates/backup.timer.j2
[Unit]
Description=Run backup for {{ backup_item.name }} daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
