# wireguard/wireguard.container
[Unit]
Description=WireGuard (wg‑easy) via slirp4netns (rootless)
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=wireguard
Image=ghcr.io/wg-easy/wg-easy:latest

# Use rootless user‑mode NAT for Internet egress
Network=pasta:\
--map-gw,\
--interface=ens5,\
--udp-ports=51820:51820,\
--tcp-ports=51821:51821

# WireGuard volume & ports
Volume=wireguard-config:/etc/wireguard:Z

# Capabilities
AddCapability=NET_ADMIN
AddCapability=NET_RAW
DropCapability=MKNOD
DropCapability=AUDIT_WRITE

# In‑container forwarding
Sysctl=net.ipv4.ip_forward=1
Sysctl=net.ipv4.conf.all.src_valid_mark=1

# WireGuard environment
Environment=WG_HOST=vpn.jackturk.dev
Environment=WG_ALLOWED_IPS=0.0.0.0/0,::/0
Environment=WG_PERSISTENT_KEEPALIVE=25
Environment=WG_DEFAULT_DNS=8.8.8.8,8.8.4.4
Environment=INSECURE=true
Environment=UI_TRAFFIC_STATS=true
Environment=UI_CHART_TYPE=2

[Install]
WantedBy=default.target