# wireguard/wireguard-ui.socket
[Unit]
Description=WireGuard UI HTTP Socket

[Socket]
ListenStream=3.15.122.25:51821

[Install]
WantedBy=sockets.target