# wireguard/wireguard.socket
[Unit]
Description=WireGuard Socket (Rootless)

[Socket]
ListenStream=[::]:51821
BindIPv6Only=both
SocketMode=0600
Service=wireguard.service

[Install]
WantedBy=sockets.target