# wireguard/wireguard-ui.socket
[Unit]
Description=WireGuard UI HTTP Socket

[Socket]
ListenStream=51821

[Install]
WantedBy=sockets.target