# wireguard/wireguard.socket
[Unit]
Description=Sockets for WireGuard

[Socket]
ListenStream=[::]:10086
ListenDatagram=[::]:51820
BindIPv6Only=both
SocketMode=0600
Service=wireguard.service

[Install]
WantedBy=sockets.target