# wireguard/wireguard-ui.socket
[Unit]
Description=WireGuard UI HTTP Socket

[Socket]
ListenStream=51821
Service=wireguard.service

[Install]
WantedBy=sockets.target