# wireguard/wireguard.socket
[Unit]
Description=WireGuard Socket (rootless)

[Socket]
ListenDatagram=[::]:51820
ListenStream=[::]:51821
BindIPv6Only=both
Multicast=false
Service=wireguard.service

[Install]
WantedBy=sockets.target